# supervisor.conf
; UNIT Platform — Supervisor configuration
; Copy to /etc/supervisor/conf.d/unit.conf on your production server.
; Run: supervisorctl reread && supervisorctl update && supervisorctl start unit-horizon

[program:unit-horizon]
command=/usr/bin/php /var/www/unit-platform/artisan horizon
directory=/var/www/unit-platform
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/unit-platform/storage/logs/horizon.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=3600
user=www-data

[program:unit-scheduler]
command=/usr/bin/php /var/www/unit-platform/artisan schedule:work
directory=/var/www/unit-platform
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/unit-platform/storage/logs/scheduler.log
user=www-data

[group:unit]
programs=unit-horizon,unit-scheduler
