# deploy/supervisor/unit-horizon.conf
;
; Supervisor config for UNIT platform — manages Laravel Horizon only.
;
; Horizon internally manages all worker pools defined in config/horizon.php.
; Do NOT add separate [program:queue-worker] entries — Horizon owns the workers.
;
; Install:
;   sudo cp deploy/supervisor/unit-horizon.conf /etc/supervisor/conf.d/unit-horizon.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start unit-horizon
;
; On deploy (after git pull + composer install + artisan migrate):
;   sudo supervisorctl restart unit-horizon
;   OR use zero-downtime: php artisan horizon:terminate  (Horizon finishes current jobs first)
;

[program:unit-horizon]
process_name    = %(program_name)s
command         = php /var/www/unit-platform/artisan horizon
directory       = /var/www/unit-platform
autostart       = true
autorestart     = true
user            = www-data
redirect_stderr = true
stdout_logfile  = /var/log/supervisor/unit-horizon.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 5
stopwaitsecs    = 3600  ; give Horizon time to drain in-flight jobs on shutdown
killasgroup     = true
stopasgroup     = true
